# .config/autostart/hyprland.desktop
[Desktop Entry]
Name=Hyprland
Comment=Hyprland
Exec=Hyprland
Type=Application
